# template_service/_supervisord.conf
[program:template]
command=bash launchService.sh template
stopasgroup=true
autostart=false
